# conf/supervisor/monitor.conf
[program:alpha_monitor]                                                                                                                                                   
command = bash /home/pi/.virtualenvs/monitor/bin/gunicorn_monitor 	    ; Command to start app
user = pi                                                 		        ; User to run as
stdout_logfile = /etc/supervisor/logs/gunicorn_supervisor_monitor.log	; Where to write log messages
redirect_stderr = true                                          	    ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                 	    ; Set UTF-8 as default encoding
